# dist/installer/templates/templates/lossless-claude-vllm-mlx.service
[Unit]
Description=lossless-claude vLLM-MLX embedding server
After=network.target

[Service]
Type=simple
ExecStart=__VLLM_MLX_BIN__ serve --embedding-model __EMBED_MODEL__ --port __MODEL_PORT__ --host 127.0.0.1
Restart=always
RestartSec=5
StandardOutput=append:__LOG_DIR__/vllm-mlx.stdout.log
StandardError=append:__LOG_DIR__/vllm-mlx.stderr.log
Environment=HOME=__HOME__

[Install]
WantedBy=default.target
